# wp360-ups-ctrl-1.0.0/usr/lib/systemd/system/wp360-ups-ctrl.service
[Unit]
Description=Controlled shutdown using UPS
After=multi-user.target

[Service]
ExecStart=/usr/sbin/wp360-ups-ctrl

[Install]
WantedBy=multi-user.target
