# wp360-ups-ctrl/DEBIAN/postinst
#!/bin/bash

if [[ "$1" = "configure" ]]
then
  dpkg-trigger update-initramfs
  systemctl enable wp360-ups-ctrl
fi
