# wp360-ups-ctrl/usr/lib/systemd/system/wp360-ups-ctrl.service
[Unit]
Description=Controlled shutdown using UPS
After=multi-user.target

[Service]
Restart=always
ExecStart=/usr/sbin/wp360-ups-ctrl

[Install]
WantedBy=multi-user.target
