# wp360-ups-ctrl/usr/lib/systemd/system/codesyscontrol.service.d/gpio24.conf
[Service]
ExecStartPre=gpioset --bias=pull-down --drive=push-pull gpiochip4 24=1
ExecStartPost=gpioset --bias=pull-down --drive=push-pull gpiochip4 24=0
