# wp360-webvisu/usr/lib/systemd/system/codesyscontrol.service.d/10-codesys-user-target.conf
[Service]
ExecStartPost=systemctl --user --machine=%USER%@ start codesys.target
ExecStopPost=systemctl --user --machine=%USER%@ stop codesys.target
