# wp360-serial-map/DEBIAN/prerm
#!/bin/bash

if [[ "$1" = "configure" ]]
then
  systemctl disable --now wp360-serial-map
fi
