# wp360-fan-ctrl/usr/lib/systemd/system/wp360-fan-ctrl.service
[Unit]
Description=Control WP360 fans
After=multi-user.target

[Service]
ExecStart=/usr/bin/wp360-fan-ctrl

[Install]
WantedBy=multi-user.target
