# wp360-webvisu-1.0.0/usr/lib/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
# The empty line is used to remove the existing ExecStart in the base service
ExecStart=
ExecStart=-/sbin/agetty --autologin user --noclear %I $TERM
